# com.thomas.mac-cleanup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.thomas.mac-cleanup</string>

  <key>ProgramArguments</key>
  <array>
    <string>~/.local/bin/mac-cleanup</string>
  </array>

  <!-- Toutes les 4 semaines, samedi à 10h00 -->
  <key>StartCalendarInterval</key>
  <dict>
    <key>Weekday</key>
    <integer>6</integer>
    <key>Hour</key>
    <integer>10</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>

  <!-- Rattrapage si la machine était éteinte -->
  <key>LaunchOnlyOnce</key>
  <false/>

  <key>StandardOutPath</key>
  <string>~/.local/share/mac-cleanup/launchd.log</string>

  <key>StandardErrorPath</key>
  <string>~/.local/share/mac-cleanup/launchd.log</string>

  <!-- Ne démarre pas si sur batterie -->
  <key>ProcessType</key>
  <string>Background</string>
</dict>
</plist>
